# debug-yolo
#!/bin/sh
### BEGIN INIT INFO
# Provides:          debug-yolo
# Required-Start:	 $remote_fs $syslog $networking
# Required-Stop:	 $remote_fs $syslog 
# Default-Start:     2 3 4 5
# Default-Stop:      
# Short-Description: Ultra minimal debugging
### END INIT INFO

# Direct output to console
exec > /dev/console 2>&1

CONTINUE_FILE=/tmp/continue-boot

log() {
    echo "DEBUG-YOLO: $1"
}


create_auto_continue() {
    # Create a file that will automatically continue boot after timeout
    log "Setting up auto-continue timer (5 minutes)"
    (
        sleep 300  # 5 minutes
        log "Timeout reached - continuing boot process automatically"
        touch $CONTINUE_FILE
    ) &
}

start() {
    log "Starting ultra-minimal debug script..."
    
    # Create the timeout handler
    create_auto_continue
    
    # Wait for continue file
    log "========================================================"
    log "BOOT PROCESS PAUSED FOR DEBUG ACCESS"
    log ""
    log "SSH ACCESS: ssh -p 10022 root@localhost"
    log ""
    log "The system will continue booting when you create:"
    log "    $CONTINUE_FILE"
    log " or run the 'continue-boot' command"
    log ""
    log "Or after 5 minutes (automatic timeout)"
    log "========================================================"
    
    # Main wait loop
    while [ ! -f $CONTINUE_FILE ]; do
        sleep 5
        log "Waiting for $CONTINUE_FILE to be created..."
    done
    
    log "Continue file found! Resuming boot process..."
}

case "$1" in
    start)
        start
        ;;
    *)
        echo "Usage: $0 {start}" >&2
        exit 1
        ;;
esac

exit 0
